# docker/laraveltest/crontab
* * * * * root /usr/local/bin/php /var/www/laraveltest/artisan schedule:run >> /dev/null 2>&1
